# qr3_82.conf
[program:qr3_82]
command=python3 /home/curator3/qr3_82/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_82/test.err.log
stdout_logfile=/home/curator3/qr3_82/test.err.log
